# root/etc/s6-overlay/s6-rc.d/cronicle-prepare/up
exec s6-setuidgid cronicle /cronicle-prepare.sh