# root/etc/s6-overlay/s6-rc.d/cronicle-prepare/up
exec s6-setuidgid mrmeeb /cronicle-prepare.sh