# root/etc/s6-overlay/s6-rc.d/cronicle/finish
echo "$e" > /run/s6-linux-init-container-results/exitcode